# dot_config/systemd/user/nix-update.timer
[Unit]
Description=Nix Updates and Maintenance
Documentation=man:nix(1)

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
