# dot_config/autostart/firewall-applet.desktop
[Desktop Entry]
Exec=/usr/bin/firewall-applet
Hidden=true
